# etc/cron.d/mackup
# /etc/cron.d/mackup - crontab for mackup(8)

0 2 * * * root /usr/share/mackup/cron/daily
